# base/further-steps/examples/dyndns/cronjob/crontab
# /etc/crontab
# DynDNS update for freedns.afraid.org
*/5 * * * * /usr/bin/curl -s https://v6.sync.afraid.org/u/YourSecretIPv6Token/
