# base/further-steps/examples/dyndns/cronjob/crontab
# Add to /etc/crontab

# Updates the dyndns entry of your domain with the correct public IP
*/5 * * * * curl -6 https://example.com/exampleValue
