# base/operating-system/linux-arch/examples/motd/motd.timer
[Unit]
Description=Updates the package database and shows the result in the motd

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
