# base/further-steps/examples/dyndns/systemd-timer/dyndns.service
# Updates the dyndns entry of your domain with the corect public IP
[Service]
ExecStart=curl -6 https://example.com/exampleValue
