# base-arch-linux/examples/dyndns/dyndns.timer
[Unit]
Description=Updates the dyndns entry of your domain with the corect public IP

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
