# base/further-steps/examples/podman/prune/systemd-timer/podman-prune.timer
#~/.config/systemd/user/podman-prune.timer
[Unit]
Description=Prune unused images/-tags to free up disk space

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
