# base/further-steps/examples/podman/prune/systemd-timer/podman-prune.service
#/etc/systemd/system/podman-prune.service
[Service]
ExecStart=/usr/bin/podman system prune -f
