# base/further-steps/examples/podman/prune/cronjob/crontab
#/etc/crontab
# Prune unused dockerimages/-tags to free up disk space
0 0 * * * /usr/bin/podman system prune -f
