# base/further-steps/examples/docker/prune/systemd-timer/docker-prune.service
# Updates the package database and shows the result in the motd

[Service]
ExecStart=/usr/bin/docker system prune -f
