# base/further-steps/examples/docker/prune/cronjob/crontab
# Add to /etc/crontab

# Prune unused dockerimages/-tags to free up disk space
0 0 * * * /usr/bin/docker system prune -f
